# deploy/systemd/tradeanalytics.service
[Unit]
Description=TradeAnalytics (Django API + Vite frontend)
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=%h/PortfolioProject

# Change these if you want different bind addresses / ports.
Environment=BACKEND_HOST=127.0.0.1
Environment=BACKEND_PORT=8000
Environment=FRONTEND_HOST=0.0.0.0
Environment=FRONTEND_PORT=5173
Environment=CUDA_VISIBLE_DEVICES=-1
Environment=TF_CPP_MIN_LOG_LEVEL=3

# Keep logs in `journalctl -u tradeanalytics -f`
ExecStart=/usr/bin/env bash %h/PortfolioProject/scripts/ubuntu/dev.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
